# tools/aur/ank-agent.service
[Unit]
Description=Ankaios agent
After=network.target ank-server.service
Wants=network.target

[Service]
Environment="RUST_LOG=info"
ExecStart=/usr/bin/ank-agent

[Install]
WantedBy=default.target
